# etc/systemd/system/letshide.service
[Unit]
Description=Let's Hide!

[Service]
Type=oneshot
EnvironmentFile=/etc/default/letshide
ExecStart=/usr/local/bin/letshide
